# omtencoder/omtencoder.service
[Unit]
Description=OMT Capture (Rust)
After=network.target

[Service]
WorkingDirectory=/opt/omtencoder
ExecStart=/opt/omtencoder/omtcapture
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
