# omtcapture/omtcapture.service
[Unit]
Description=OMT Capture
After=network.target

[Service]
ExecStart=/opt/omtcapture/omtcapture

[Install]
WantedBy=multi-user.target
